# ecosystem-stack/scripts/local-ecosystem-cpu-temp.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- Legacy reference only. Prefer: bash ecosystem-stack/scripts/macos-host-metrics-scheduler.sh install
     (wired from ecosystem-stack/services/dashboard.sh on macOS). That generates
     ~/Library/LaunchAgents/com.local-ecosystem.host-cpu-temp.plist with correct paths. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local-ecosystem.cpu-temp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>REPO_ROOT/ecosystem-stack/scripts/macos-write-cpu-temp.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/local-ecosystem-cpu-temp.out.log</string>
  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/local-ecosystem-cpu-temp.err.log</string>
</dict>
</plist>
